# crontab
0 6,12,18,0 * * * cd /opt/necpProd;scrapy runspider myspider_gfs.py >> /var/log/myspider_gfs.log 2>&1
0 7 * * * cd /opt/necpProd;scrapy runspider myspider_nsst.py >> /var/log/myspider_nsst.log 2>&1